# server/xinetd.d/ldminfod
# default: off
# description: LDM Info Daemon
service ldminfod
{
	disable		= yes
	type		= UNLISTED
	socket_type     = stream
	wait            = no
	user            = nobody
	server          = /usr/sbin/ldminfod
	port		= 9571
}
