# server/xinetd.d/nbdrootd
# default: off
# description: NBD root daemon
service nbdrootd
{
	disable		= yes
	type		= UNLISTED
	flags		= KEEPALIVE
	socket_type     = stream
	wait            = no
	user            = nobody
	server          = /usr/sbin/nbdrootd
	port		= 2000
}
