# server/xinetd.d/nbdswapd
# default: off
# description: NBD swap daemon
service nbdswapd
{
	disable		= yes
	type		= UNLISTED
	socket_type     = stream
	wait            = no
	user            = nobody
	server          = /usr/sbin/nbdswapd
	port		= 9572
}
