# client/initscripts/Gentoo/ltsp-client.initd
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
    need ltsp-client-setup
    after alsasound
}

start() {
	ebegin "Starting LTSP Client"
	. /usr/share/ltsp/ltsp-init-common

	start_screen_sessions || true
	start_sound || true
}
